# manual-kubernetes-on-arm/modules/bootkube/modules/control-plane/modules/cloud-init/files/etc/systemd/system/kubelet.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/kubelet \
    --container-runtime=remote \
    --container-runtime-endpoint=unix:///run/containerd/containerd.sock \
    --logging-format=text

#    --bootstrap-kubeconfig
#    --config string \
#    --kubeconfig string \
