# manual-kubernetes-on-arm/modules/kubernetes-worker/cloud-init/files/etc/systemd/system/kubelet.service.d/10-crio.conf
[Unit]
After=crio.service
Requires=crio.service
