# challenge/ansible/roles/kubernetes/templates/etc/systemd/system/kubelet.service.d/override.conf.j2
[Unit]
After=crio.service
Requires=crio.service

[Service]
ExecStart=
ExecStart=/usr/bin/kubelet \
  --cloud-provider=aws \
  --config=/etc/kubernetes/kubelet-config.yaml \
  --container-runtime=remote \
  --container-runtime-endpoint=unix:///var/run/crio/crio.sock \
  --kubeconfig=/etc/kubernetes/kubelet.config \
  --network-plugin=cni \
  --register-node=true \
  --v=2

Restart=on-failure
RestartSec=5
